# rules.d/90-acq2106.rules
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07H", SYMLINK+="tty_acq2106_201" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07J", SYMLINK+="tty_acq2106_203" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07L", SYMLINK+="tty_acq2106_204" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07M", SYMLINK+="tty_acq2106_205" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07M", SYMLINK+="tty_acq2106_205" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07N", SYMLINK+="tty_acq2106_206" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07O", SYMLINK+="tty_acq2106_207" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07G", SYMLINK+="tty_acq2106_208" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AH07D07K", SYMLINK+="tty_acq2106_209" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01D703", SYMLINK+="tty_acq2106_210" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01BN99", SYMLINK+="tty_acq2106_211" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01BOJ3", SYMLINK+="tty_acq2106_212" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01DC5Z", SYMLINK+="tty_acq2106_213" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01AREE", SYMLINK+="tty_acq2106_214" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01B8NU", SYMLINK+="tty_acq2106_215" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01AMKY", SYMLINK+="tty_acq2106_216" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01AEPZ", SYMLINK+="tty_acq2106_217" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01DNNE", SYMLINK+="tty_acq2106_218" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01BNW8", SYMLINK+="tty_acq2106_219" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01BOPJ", SYMLINK+="tty_acq2106_220" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01BOJP", SYMLINK+="tty_acq2106_221" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01DMPQ", SYMLINK+="tty_acq2106_222" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01DBPU", SYMLINK+="tty_acq2106_223" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01DDNX", SYMLINK+="tty_acq2106_224" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01DEDX", SYMLINK+="tty_acq2106_225" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01BOKR", SYMLINK+="tty_acq2106_226" RUN+="/usr/local/bin/make-tty-symlink"
SUBSYSTEMS=="usb", ATTRS{serial}=="AU01DTX1", SYMLINK+="tty_acq2106_227" RUN+="/usr/local/bin/make-tty-symlink"

