# rules.d/90-acq2106.rules
SUBSYSTEMS=="usb", ATTRS{serial}=="AM01L51H", SYMLINK+="tty_acq2106_001" RUN+="/usr/local/bin/make-tty-symlink"
