# rules.d/90-acq1001.rules
SUBSYSTEMS=="usb", ATTRS{serial}=="AM01L50P", SYMLINK+="tty_acq1001_001" RUN+="/usr/local/bin/make-tty-symlink"
