# rules.d/90-cpsc2.rules
SUBSYSTEMS=="usb", ATTRS{serial}=="A107W5MP", SYMLINK+="tty_cpsc2_004" RUN+="/usr/local/bin/make-tty-symlink"
